# gunicorn.conf
[program:gunicorn]
command=/usr/bin/gunicorn wsgi:application -b localhost:5000 -c gunicorn.conf
directory=/deploy/app
